# scripts/claude-ideas/claude-build-watcher.service
[Unit]
# Install as a *user* service (~/.config/systemd/user/) so it can spawn
# per-build transient units via `systemd-run --user`. Requires linger:
#     sudo loginctl enable-linger luis
Description=Claude Build Watcher — picks up dashboard build jobs and runs Claude Code
After=default.target

[Service]
Type=simple
ExecStart=/opt/homelab/scripts/claude-ideas/build-watcher.sh
Restart=always
RestartSec=10
StandardOutput=append:/opt/claude-ideas/build-watcher.log
StandardError=append:/opt/claude-ideas/build-watcher.log

[Install]
WantedBy=default.target
